# packer/files/puma.service
[Unit]
Description=puma

[Service]
WorkingDirectory=/home/baskakov/reddit
ExecStart=/usr/local/bin/puma

[Install]
WantedBy=multi-user.target
